# bot/agent-bot.service
[Unit]
Description=Agent Bot - Personal AI Telegram Bot
After=network.target

[Service]
Type=simple
User=agent
WorkingDirectory=/home/agent/.agent/bot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
EnvironmentFile=/home/agent/.agent/.env

[Install]
WantedBy=multi-user.target
